# configs/examples/ferrumd.service
[Unit]
Description=FerrumGate daemon
Documentation=https://github.com/ferrumgate/ferrum-gate/blob/main/docs/guides/operator.md
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ferrumgate
Group=ferrumgate
EnvironmentFile=-/etc/ferrumgate/ferrumd.env
ExecStart=/usr/local/bin/ferrumd --config /etc/ferrumgate/ferrumgate.toml
Restart=on-failure
RestartSec=5s
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/var/lib/ferrumgate /var/log/ferrumgate

[Install]
WantedBy=multi-user.target
