# configs/examples/postgres-backup.timer
[Unit]
Description=FerrumGate PostgreSQL backup timer (15-minute interval)
Documentation=./docs/implementation-path/109-p5c-postgresql-backup-restore-runbook.md

[Timer]
OnCalendar=*:0/15
Persistent=true
RandomizedDelaySec=60

[Install]
WantedBy=timers.target
